# schedules/crontab
PATH=/usr/local/sbin:/usr/local/bin:/home/user1/bin
MAILTO=dmitri@vprex.com
*/50 * * * * /usr/bin/php /home/mainstreet/yiicmd.php cron null
